# 2_task/bagcleaner.service
[Unit]
Description=Script to clean old bag files

[Service]
Type=oneshot
ExecStart=/path/to/your/script.sh

[Install]
WantedBy=multi-user.target
